# docker/php/crontab/laravel
* * * * * root /usr/local/bin/php /work/artisan schedule:run >> /dev/null 2>&1
